# net-fs/netatalk/files/cnid_metad.init.2
#!/sbin/runscript
# Copyright 1999-2011 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header$

: ${CNID_CONFIG:=-l log_note}

command="/usr/sbin/${SVCNAME}"
command_args="${CNID_CONFIG}"
pidfile="/var/run/${SVCNAME}.pid"

depend() {
	need net
    use logger
}
