# net-fs/netatalk/files/netatalk.init
#!/sbin/runscript
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header$

command="/usr/sbin/${SVCNAME}"

depend() {
	need net
	need avahi-daemon
	use logger dns
	after entropy
}
